# etc/systemd/system/monitor-sh.service
[Unit]
Description="Monitoring the network"
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor.sh